# consumer/kafka-consumer.conf
[program:kafka-consumer]
command=php /app/consumer.php
autostart=true
autorestart=true
stderr_logfile=/var/log/kafka-consumer.err.log
stdout_logfile=/var/log/kafka-consumer.out.log
